# init/watcher.conf
description "Simple CPU and Memory usage tracking"
author "brosbunbo"

start on filesystem or runlevel [2345]
stop on shutdown

script
    exec /root/Watcher/start
end script